# train-poller.service
[Unit]
Description=Telegram Train Subscription Poller
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/railbot/telegram-train-tracker/
ExecStart=/home/railbot/telegram-train-tracker/venv/bin/python /home/railbot/telegram-train-tracker/run_poller.py --daemon
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Ensure environment variables are loaded
EnvironmentFile=/home/railbot/telegram-train-tracker/.env

[Install]
WantedBy=multi-user.target
