# train-poller.service
[Unit]
Description=Telegram Train Subscription Poller
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/railbot/telegram-train-tracker/
ExecStart=/usr/bin/python3 /home/railbot/telegram-train-tracker/run_poller.py --daemon --interval 300
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Ensure environment variables are loaded
EnvironmentFile=/Volumes/workplace/cline-rail-bot/.env

[Install]
WantedBy=multi-user.target
