# train-bot.service
[Unit]
Description=Telegram Train Information Bot
After=network.target

[Service]
Type=simple
User=itaiaz
WorkingDirectory=/Volumes/workplace/cline-rail-bot
Environment=PYTHONPATH=/Volumes/workplace/cline-rail-bot
ExecStart=/usr/bin/python3 /Volumes/workplace/cline-rail-bot/run.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
